# deploy/systemd/backup-run.timer
# Timer diario do backup (doc 9: timer substitui o cron). Ajuste OnCalendar por cliente.
[Unit]
Description=Dispara o backup diario

[Timer]
OnCalendar=*-*-* 22:00:00
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
